# raspberrypi2z/net-watchdog/net-watchdog.service
[Unit]
Description=Network watchdog - reboot if router unreachable
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/net-watchdog.sh
